# roles/setup_docker_haproxy/files/update_haproxy_config.conf
[program:update_haproxy_config]
command=/home/ubuntu/update_haproxy_config.py
