# roles/setup_check_my_ip_server/templates/check_my_ip_super_config.conf
[program:check_my_ip]
command=/home/ubuntu/check_my_ip_server.py
